# run-scripts/crontab
*/30 * * * * /home/opensim/checkOldLogFiles.sh

