# opensim-region-old/crontab
*/2 * * * * /home/opensim/checkOpenSim.sh
*/30 * * * * /home/opensim/checkOldLogFiles.sh

